# celery-beat.ini
[program:domainer-celery-beat]
command=celery -A domainer beat -l info --scheduler django_celery_beat.schedulers:DatabaseScheduler
process_name=domainer-celery-beat-%(process_num)d
numprocs=1
directory=/opt/wwwroot/ops-domainer-api/
autostart=true
autorestart=true
autoretrice=10
user=publish
environment=ENV="prod"
numprocs=1
redirect_stderr=true
stdout_logfile=/opt/logs/ops-domainer-api/celery-beat.log
stopwatisecs=60
priority=994
stdout_logfile_maxbytes = 50MB