# ops-domainer-api.ini
[program:ops-domainer-api]
command=/data/python-evev/ops-domainer-api/bin/python3.11 manage.py runserver 0.0.0.0:8080
process_name=%(program_name)s-%(process_num)d
directory=/data/wwwroot/ops-domainer-api/
autostart=true
autorestart=true
startsecs=10
startretries=3
exitcodes=0,2
environment=ENV="prod"
user=wwwuser
numprocs=1 ; 定义1个进程
redirect_stderr=true
stdout_logfile=/data/logs/ops-domainer-api/ops-domainer-api.log
stopwaitsecs=60
priority=994
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
[program:ops-domainer-api-celery-beat]
command=/data/python-evev/ops-domainer-api/bin/celery -A domainer beat -l info --scheduler django_celery_beat.schedulers:DatabaseScheduler
process_name=%(program_name)s-celery-beat-%(process_num)d
directory=/data/wwwroot/ops-domainer-api/
autostart=true
autorestart=true
startsecs=10
startretries=3
exitcodes=0,2
environment=ENV="prod"
user=wwwuser
numprocs=1 ; 定义1个进程
redirect_stderr=true
stdout_logfile=/data/logs/ops-domainer-api/celery-beat.log
stopwaitsecs=60
priority=994
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
[program:ops-domainer-api-celery-worker]
command=/data/python-evev/ops-domainer-api/bin/celery -A domainer worker -l info
process_name=%(program_name)s-celery-worker-%(process_num)d
directory=/data/wwwroot/ops-domainer-api/
autostart=true
autorestart=true
startsecs=10
startretries=3
exitcodes=0,2
environment=ENV="prod"
user=wwwuser
numprocs=2 ; 定义3个进程
redirect_stderr=true
stdout_logfile=/data/logs/ops-domainer-api/celery-worker.log
stopwaitsecs=60
priority=994
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10